# container_configs/nextcloud/etc/systemd/system/nginx.service.d/chroot.conf
[Service]
#ProtectSystem=full
ReadWritePaths=/var/run/php/ /var/log/nginx/ /var/www/nextcloud/custom_apps/ /var/www/nextcloud/data/
ReadOnlyPaths=/usr/sbin/nginx /var/www/ /etc/nginx
PrivateTmp=true
PrivateDevices=true
ProtectHome=false
ProtectKernelTunables=true
NoNewPrivileges=true
ProtectKernelModules=true

