# container_configs/nextcloud/etc/systemd/system/nextcloud-cron.timer
[Unit]
Description=Run nextcloud-cron every 4 minutes

[Timer]
OnBootSec=4min
OnUnitActiveSec=4min

[Install]
WantedBy=multi-user.target
